# packaging/linux/ziplock-backend.service
[Unit]
Description=ZipLock Password Manager Backend Service
Documentation=https://github.com/ejangi/ziplock
After=network.target
Wants=network.target

[Service]
Type=notify
User=ziplock
Group=ziplock
ExecStart=/usr/bin/ziplock-backend --daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true

# File system access
ReadWritePaths=/var/lib/ziplock
ConfigurationDirectory=ziplock
StateDirectory=ziplock
LogsDirectory=ziplock
CacheDirectory=ziplock

# Environment
Environment=RUST_LOG=info
Environment=ZIPLOCK_CONFIG_DIR=/etc/ziplock

# Capabilities
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
